# root/etc/etc-sub/skel/.config/autostart/steam-or-setup.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Steam Big Picture
Icon=/usr/share/icons/steambigpicture.png
Exec=sh -c "/etc/skel/steam-or-setup.sh" 
Categories=Network;Multimedia;Game;
